# config/etc/services.d/ssdb/run
#!/bin/sh

/usr/local/ssdb/ssdb-server /etc/sysconfig/ssdb.conf